# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/redis.err.log
stdout_logfile=/var/log/redis.out.log
user=appuser

[program:trading-system]
command=python /app/pipeline.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/system.err.log
stdout_logfile=/app/logs/system.out.log
user=appuser
environment=PYTHONPATH="/app"

[program:api-server]
command=python /app/api_server.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/api.err.log
stdout_logfile=/app/logs/api.out.log
user=appuser
environment=PYTHONPATH="/app"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
